# acserver.cron
* * * * * /usr/bin/php /var/www/acserver/index.php batch_update_from_carddb 2>&1 >> /var/log/acserver_update.log
